# packaging/systemd/cbp-dashboard.service
# CryptKeep operator dashboard (read-only visibility; bind per your host).
[Unit]
Description=CryptKeep dashboard
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=exec
User=cbp
Group=cbp
WorkingDirectory=/opt/crypto-bot-pro
EnvironmentFile=/etc/cbp/cbp.env
ExecStart=/opt/crypto-bot-pro/.venv/bin/python scripts/run_dashboard.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/cbp

[Install]
WantedBy=multi-user.target
